# raspberry-pi/systemd/focus-timer.service
[Unit]
Description=Focus Timer Servo Controller
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/focus-timer-iot/raspberry-pi
ExecStart=/usr/bin/python3 /home/pi/focus-timer-iot/raspberry-pi/focus_servo.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
